# systemd/atistik-raspberry-analyze-guard.service
[Unit]
Description=Verify Atistik Raspberry morning analysis and retry locally
After=atistik-raspberry-api.service network-online.target
Wants=atistik-raspberry-api.service network-online.target

[Service]
Type=oneshot
User=pi
Group=pi
WorkingDirectory=/opt/atistik/backend
ExecStart=/opt/atistik/backend/scripts/raspberry/ensure-automation.sh analyze
TimeoutStartSec=7200
