# systemd/atistik-raspberry-analyze.service
[Unit]
Description=Atistik Raspberry morning analysis
After=atistik-raspberry-api.service network-online.target
Wants=atistik-raspberry-api.service network-online.target

[Service]
Type=oneshot
User=pi
Group=pi
WorkingDirectory=/opt/atistik/backend
ExecStart=/opt/atistik/backend/scripts/raspberry/run-automation.sh analyze
TimeoutStartSec=7200
