# systemd/atistik-raspberry-api.service
[Unit]
Description=Atistik Raspberry API container
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/atistik/backend
ExecStart=/usr/bin/docker compose -f docker-compose.raspberry.yml up -d atistik-api
ExecStop=/usr/bin/docker compose -f docker-compose.raspberry.yml stop atistik-api
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
