# systemd/atistik-raspberry-results.service
[Unit]
Description=Atistik Raspberry evening results submission
After=atistik-raspberry-api.service network-online.target
Wants=atistik-raspberry-api.service network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/atistik/backend
ExecStart=/opt/atistik/backend/scripts/raspberry/run-automation.sh results
TimeoutStartSec=7200
